# systemd/room-signage.service
[Unit]
Description=Room Signage Service
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
# Run as the roomsign user which owns the application files.
User=roomsign
Group=roomsign
WorkingDirectory=/opt/room-signage/app

# Load environment variables from the configuration file. This file should
# contain lines like GOOGLE_SERVICE_ACCOUNT_JSON=/opt/room-signage/config/service-account.json
EnvironmentFile=/opt/room-signage/config/room-signage.env

ExecStart=/opt/room-signage/venv/bin/uvicorn app.main:app --host 127.0.0.1 --port 8080

# Restart the service if it fails.
Restart=always
RestartSec=5s

# Security hardening. Limit the service’s view of the system.
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
# Allow the service to read and write its own logs directory.
ReadWritePaths=/opt/room-signage/logs

[Install]
WantedBy=multi-user.target